# deploy/wind-stats.service
[Unit]
Description=Scrape windguru
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/bin/wind-stats.sh
